# deploy/dongbo-tiktok-monitor.timer
[Unit]
Description=Run the Malaysia TikTok public product monitor regularly

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
RandomizedDelaySec=2min
Persistent=true
Unit=dongbo-tiktok-monitor.service

[Install]
WantedBy=timers.target
