# deploy/dongbo-tiktok-monitor.service
[Unit]
Description=Dongbo ERP Malaysia TikTok public product monitor
After=network-online.target dongbo-erp.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/dongbo/app/backend
EnvironmentFile=/opt/dongbo/app/.env
ExecStart=/opt/dongbo/venv/bin/python manage.py monitor_tiktok_products --market MY
NoNewPrivileges=true
PrivateTmp=true
